# so-kibana/files/kibana.service
[Unit]
Description=Kibana

[Service]
Type=simple
User=kibana
ExecStart=/usr/local/bin/kibana-docker

[Install]
WantedBy=multi-user.target
